# scripts/sudo-ai-v5.service
[Unit]
Description=SUDO-AI v5 Agent Stack
After=network.target

[Service]
Type=forking
User=root
Group=root
# __SUDO_AI_HOME__ is a placeholder substituted with the real install path by
# install.sh (systemd does NOT expand ${VAR:-default}). To install manually, run:
#   sed "s|__SUDO_AI_HOME__|/your/path|g" scripts/sudo-ai-v5.service | sudo tee /etc/systemd/system/sudo-ai-v5.service
Environment=SUDO_AI_HOME=__SUDO_AI_HOME__
Environment=HOME=__SUDO_AI_HOME__
Environment=PM2_HOME=__SUDO_AI_HOME__/.pm2

WorkingDirectory=__SUDO_AI_HOME__

# Pre-start cleanup: kill any leaked tsx or esbuild orphans (from install.sh / Wave1 bootstrap too)
ExecStartPre=-/bin/sh -c "pkill -9 -f 'tsx src/cli\.ts' || true"
ExecStartPre=-/bin/sh -c "pkill -9 -f 'esbuild' || true"
ExecStartPre=-/bin/sh -c "sleep 1"
ExecStartPre=-/bin/sh -c "pm2 delete sudo-ai-v5-staging || true"
ExecStartPre=-/bin/sh -c "pm2 delete sudo-ai-v5 || true"
ExecStartPre=-/bin/sh -c "pm2 save --force || true"
ExecStartPre=-/bin/sh -c "sleep 1"

# Use SUDO_AI_HOME for ecosystem (portable; install.sh copies/enables service optionally)
ExecStart=/usr/bin/pm2 start __SUDO_AI_HOME__/ecosystem.config.cjs --only sudo-ai-v5 --update-env
ExecStop=/usr/bin/pm2 stop sudo-ai-v5
ExecReload=/usr/bin/pm2 reload sudo-ai-v5

# If the service dies, systemd restarts it (pm2 itself also restarts, but this
# is the ultimate safety net).
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
